# skill/research-assist/automation/launchd/com.zephyr.research-assist.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zephyr.research-assist.digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/zephyr/Bio/Codex/skill/research-assist/automation/launchd/run_daily_digest.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zephyr/Bio/Codex/skill/research-assist</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zephyr/Bio/Codex/skill/research-assist/automation/logs/daily-digest.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zephyr/Bio/Codex/skill/research-assist/automation/logs/daily-digest.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
